# deploy/gunicorn.service
[Unit]
Description=Gunicorn daemon for VinFlow Django application
After=network.target

[Service]
Type=notify
# User and group to run as (change to your user)
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/var/www/vinflow

# Environment variables
Environment="PATH=/var/www/vinflow/venv/bin"
EnvironmentFile=/var/www/vinflow/.env

# Gunicorn command
ExecStart=/var/www/vinflow/venv/bin/gunicorn \
    --config /var/www/vinflow/gunicorn_config.py \
    core.wsgi:application

# Restart policy
Restart=always
RestartSec=10

# Process limits
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vinflow-gunicorn

# Create runtime directory for PID file
RuntimeDirectory=vinflow
RuntimeDirectoryMode=0755

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

